# crontab
30 * * * * root /usr/local/bin/python3 /app/pure_adaptor/pure_adaptor.py > /proc/1/fd/1 2>/proc/1/fd/2
